# pocketchip-batt.service
[Unit]
Description=Handle battery, backlight and PocketCHIP keymap

[Service]
Type=oneshot
ExecStart=/usr/sbin/pocketchip-batt daemon
